# docker/laravel-env_php8/config/supervisord/nginx.ini
[program:nginx_supervisord]
command=/usr/sbin/nginx -g "daemon off;"
stderr_logfile=/var/log/nginx_supervisord.err.log
stdout_logfile=/var/log/nginx_supervisord.out.log