# docker/laravel-env/php8/config/supervisord/php-fpm.ini
[program:php_supervisord]
command=/usr/sbin/php-fpm -RF
stderr_logfile=/var/log/php_supervisord.err.log
stdout_logfile=/var/log/php_supervisord.out.log